# systemd/track-setup.service
[Unit]
Description=TRACK Runtime Setup (tmpfs dirs and config)
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '\
    mkdir -p /tmp/track-logs && \
    cp -n /opt/track/config/display.dbc /tmp/display.dbc 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
