# systemd/track-cloud-bridge.service
[Unit]
Description=TRACK Cloud Bridge
After=network-online.target track-can-reader.service
Wants=network-online.target
Requires=track-can-reader.service

[Service]
Type=simple
ExecStart=/opt/track/cloud-bridge
Environment=CB_URL=wss://track-web.fly.dev/ws/pi
Environment=CB_DEVICE_ID=dev-001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
